# trinity-acp.service
[Unit]
Description=Trinity ACP Agent API Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/acp-gridcore
Environment="PATH=/home/ubuntu/acp-gridcore/venv/bin"
ExecStart=/home/ubuntu/acp-gridcore/venv/bin/uvicorn api_server:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/acp-gridcore/api_server.log
StandardError=append:/home/ubuntu/acp-gridcore/api_server.log

[Install]
WantedBy=multi-user.target
